# etc/systemd/system/bind9-blocklists-update.timer
[Unit]
Description=Bind9 Blocklist update timer - Daily

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
