# etc/systemd/system/bind9-blocklists-update.service
[Unit]
Description=Bind9 Block list update service
After=network.target

[Service]
ExecStart=/opt/bind9/update-blocklists.sh
