# etc/systemd/system/unbound-blocklist-update.service
[Unit]
Description=Unbound blocklist update service
After=network.target

[Service]
ExecStart=/opt/unbound/blocklists/update.sh
