# etc/systemd/system/unbound-root-hints-update.timer
[Unit]
Description=Unbound roothints update timer - Weekly

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
