# etc/systemd/system/unbound-root-hints-update.service
[Unit]
Description=Unbound root hints update service
After=network.target

[Service]
ExecStart=/opt/unbound/update-roothints.sh
