# etc/systemd/system/unbound-root-hints-update.timer
[Unit]
Description=Unbound roothints update timer - Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
